# systemd/gradtemp.service
[Unit]
Description=Change color temperature based on current time

[Service]
Type=oneshot
ExecStart=/path/to/gradtemp
